# etc_server/profile.d/kiosk-parameters.sh
#!/bin/sh
export CHARGE_CONTROLLER=/usr/local/bin/phocos.py
#export CHARGE_CONTROLLER=/usr/local/bin/ts_45
#interface to access the internet
export NET=eth0
